# sync-gnome-theme-to-qt/usr/lib/systemd/user/sync-gnome-theme-to-qt.path
[Unit]
Description=Sync gtk theme from Gnome to QT
ConditionUser=!root

[Path]
PathModified=%h/.config/dconf/user
Unit=sync-gnome-theme-to-qt.service

[Install]
WantedBy=default.target
